# files/crontab
*/1 * * * * /usr/bin/perl /home/drew/workspace/battery/batteryMonitor.pl
