# files/crontab
*/1 * * * * perl /home/drew/workspace/battery/autoHibernate.pl
